# hardware/cpu/microcode-ctl/files/89-microcode.rules
KERNEL=="cpu[0-9]*", ACTION=="add", RUN+="/sbin/modprobe microcode"
KERNEL=="microcode", ACTION=="add", RUN+="/sbin/microcode_ctl -Qu"
