# hardware/printer/cups/files/cups-libusb.rules
# Allow the lp group to read and write raw USB devices for printers.
# This is required for the CUPS usb backend when built with libusb
# support.
ACTION=="add", SUBSYSTEM=="usb", \
  ATTR{bInterfaceClass}=="07", ATTR{bInterfaceSubClass}=="01", \
  RUN+="/usr/bin/setfacl -m g:lp:rw /dev/%P"
