# util/misc/mlocate/files/mlocate.cron
#!/bin/sh
nodevs=$(awk '$1 == "nodev" { print $2 }' </proc/filesystems)
nice --adjustment=+19 ionice -c3 /usr/bin/updatedb -f "$nodevs"
