# kernel/drivers/module-virtualbox-guest/files/48-vbox-guest.rules
KERNEL=="vboxguest", NAME="vboxguest", MODE="0660", OWNER="vboxadd", RUN+="/lib/udev/vbox-guest.sh"
KERNEL=="vboxuser", NAME="vboxuser", MODE="0666", OWNER="vboxadd"
